# etc/supervisor-zookeeper.conf
[program:zookeeper]
environment=ZOO_LOG_DIR="/var/log/zookeeper",ZOO_LOG4J_PROP="INFO,ROLLINGFILE"
command=/opt/zookeeper/bin/zkServer.sh start-foreground
autostart=true
autorestart=true
